# supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
directory=/etc/nginx
